# docker/start.service
[Unit]
Description=Customize the  script
After=network.target

[Service]
ExecStart=/start.sh

[Install]
WantedBy=multi-user.target